# pop/mtls/crontab
# m h  dom mon dow   command
0 0 1,15 * * /usr/local/moon/pop/mtls/backupdb.sh > /tmp/backupdb.log
0 1 * * * /usr/local/moon/pop/mtls/stat.sh > /tmp/stat.log
