# pop/mtls/crontab
# m h  dom mon dow   command
55 23 * * * /usr/local/moon/pop/mtls/stat.sh > /tmp/stat.log
0 0 1,15 * * /usr/local/moon/pop/mtls/backupdb.sh > /tmp/backupdb.log
